# scripts/systemd/oaao-research-cron.timer
[Unit]
Description=OAAO Article Research cron timer (every 5 minutes)
Requires=oaao-research-cron.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
